# systemd/minitelnet-client.service
[Unit]
Description=MinitelNet client
After=network-online.target
Wants=network-online.target

[Service]
User=pi
WorkingDirectory=/home/pi/minitelnet/pi-client
# All per-device config (server URL, token, apps, flags) lives in
# /etc/minitelnet/device.json — written once at enroll time. This unit is
# identical on every Pi.
ExecStart=/home/pi/minitelnet/pi-client/.venv/bin/python /home/pi/minitelnet/pi-client/launcher.py
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
